# .local/share/s6/sv/emacs-daemon/run
#!/bin/execlineb -P
exec /usr/bin/emacs --fg-daemon
